# templates/rc.local
#!/bin/sh -e

bash /var/www/ghost-serverless/rc.local
